# private_dot_local/private_share/restic-backups/private_restic-r2-system-backup.timer
[Unit]
Description=Run Restic R2 system backup periodically

[Timer]
# Run daily at 03:30 when possible.
OnCalendar=*-*-* 03:30:00
# If the machine was off/asleep at that time, run after boot.
Persistent=true
# Add jitter so it does not always run at the exact same second.
RandomizedDelaySec=30m
Unit=restic-r2-system-backup.service

[Install]
WantedBy=timers.target
